# lighttpd/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/lighttpd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/lighttpd
drwxr-xr-x	root/root	etc/ssl/
drwxr-xr-x	root/root	etc/ssl/certs/
-rw-------	root/root	etc/ssl/certs/lighttpd.pem (EMPTY)
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/lighttpd/
-rw-r--r--	root/root	usr/lib/lighttpd/mod_access.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_access.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_access.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_accesslog.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_accesslog.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_accesslog.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_ajp13.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_ajp13.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_ajp13.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_alias.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_alias.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_alias.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_auth.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_auth.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_auth.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_authn_file.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_file.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_file.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_cgi.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_cgi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_cgi.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_deflate.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_deflate.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_deflate.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_dirlisting.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_dirlisting.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_dirlisting.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_evhost.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_evhost.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_evhost.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_expire.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_expire.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_expire.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_extforward.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_extforward.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_extforward.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_fastcgi.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_fastcgi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_fastcgi.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_indexfile.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_indexfile.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_indexfile.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_openssl.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_openssl.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_openssl.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_proxy.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_proxy.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_proxy.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_redirect.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_redirect.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_redirect.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_rewrite.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rewrite.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rewrite.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_rrdtool.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rrdtool.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rrdtool.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_scgi.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_scgi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_scgi.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_setenv.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_setenv.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_setenv.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_simple_vhost.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_simple_vhost.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_simple_vhost.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_sockproxy.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_sockproxy.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_sockproxy.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_ssi.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_ssi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_ssi.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_staticfile.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_staticfile.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_staticfile.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_status.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_status.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_status.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_userdir.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_userdir.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_userdir.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_vhostdb.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_webdav.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_webdav.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_webdav.so
-rw-r--r--	root/root	usr/lib/lighttpd/mod_wstunnel.a
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_wstunnel.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_wstunnel.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/lighttpd
-rwxr-xr-x	root/root	usr/sbin/lighttpd-angel
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/lighttpd-angel.8.gz
-rw-r--r--	root/root	usr/share/man/man8/lighttpd.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/www/
drwxr-xr-x	root/root	var/www/lighttpd/
drwxr-xr-x	www/www	var/www/lighttpd/dev/
drwxr-xr-x	www/www	var/www/lighttpd/htdocs/
drwxr-xr-x	www/www	var/www/lighttpd/logs/
drwxr-xr-x	www/www	var/www/lighttpd/tmp/
